# config/supervisord/blog_server.ini
[program:blog-server]
directory = /home/a624042309/zxins
environment=PATH="/home/a624042309/zxins/ENV/bin"
command = /home/a624042309/zxins/ENV/bin/python blog_server.py
process_name = blog-server-%(process_num)s
autorestart = unexpected
redirect_stderr = true
stdout_logfile = /home/a624042309/zxins/log/blog-server-stdout.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=50
stdout_capture_maxbytes=10MB
stderr_logfile = /home/a624042309/zxins/log/blog-server-stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=50
stderr_capture_maxbytes=10MB
startsecs=5
stopwaitsecs=30